# apkovl/etc/local.d/docker.start
#!/bin/sh

exec > /dev/console
exec 2>&1

echo "Schedule docker for install ..."
echo "docker" >> /etc/apk/world
