# apkovl/etc/local.d/docker.start
#!/bin/sh

exec > /dev/console
exec 2>&1

echo "Add docker to base system"
echo "docker" >> /etc/apk/world
