# apkovl/etc/local.d/setup-sh.start
#!/bin/sh

set -e

exec > /dev/console
exec 2>&1

CMD="/root/$(basename ${0} .start)" 
[ -x "${CMD}" ] && "${CMD}"

rm -f "${0}"
