# apkovl/etc/local.d/setup-alpine.start
#!/bin/sh

exec > /dev/console
exec 2>&1

rm -f "${0}"
/root/setup-alpine.start
