# apkovl/etc/local.d/setup-alpine.start
#!/bin/sh

set -e

exec > /dev/console
exec 2>&1

rm -f "${0}"

CMD="/root/$(basename ${0} .start)" 
[ -x "${CMD}" ] && "${CMD}"
