# apkovl/etc/local.d/setup-subids.start
#!/bin/sh

/usr/local/sbin/setup-subids
