# apkovl/etc/local.d/tmp-cleaner.stop
#!/bin/sh

exec > /dev/console
exec 2>&1

find /tmp -mindepth 1 -delete
