# apkovl/etc/local.d/tmp-cleaner.stop
#!/bin/sh

find /tmp -mindepth 1 -delete
