# apkovl/etc/local.d/tmp-cleaner.start
#!/bin/sh

find /tmp -mindepth 1 -mmin +1 -delete
