# apkovl/pxe/etc/init.d/setup-alpine
#!/sbin/openrc-run

depend() {
	need net
	after firewall
}

start() {
	setup-alpine -q -f /root/setup-alpine.answers
}
