# apkovl/pxe/etc/init.d/setup-alpine
#!/sbin/openrc-run

depend() {
	keyword -stop -timeout
}

start() {
	ebegin "Running setup-alpine $ALPINEOPTS"
	/sbin/setup-alpine $ALPINEOPTS
	eend $?
}
