# apkovl/pxe/etc/local.d/hello-world.start
#!/bin/sh
echo "Hello normal world!"
echo "Hello console world!" > /dev/console
echo "Hello stdout world!" > /dev/stdout
