# apkovl/pxe/etc/local.d/setup-alpine.start
#!/bin/sh
setup-alpine -q -f /root/setup-alpine.answers
