# apkovl/pxe/etc/local.d/setup-alpine.start
#!/bin/sh
setup-alpine -e -f /root/setup-alpine.answers
