# apkovl/pxe/etc/local.d/setup-alpine.start
#!/bin/sh
echo y | setup-alpine -e -f /root/setup-alpine.answers
