# apkovl/etc/periodic/daily/apk
#!/bin/sh

set -eu
sleep $(( $RANDOM % 7200 ))

apk -U upgrade
