# apkovl/etc/periodic/daily/apk
#!/bin/sh

set -eu

exec apk -U upgrade
