# apkovl/etc/periodic/daily/apk
#!/bin/sh

set -eu

sleep $(expr $RANDOM % 7200)
exec apk -U upgrade
