# apkovl/etc/periodic/daily/gem
#!/bin/sh

set -eu

sleep $(expr $RANDOM % 7200)

gem update --system
gem update
gem cleanup
