# apkovl/etc/periodic/daily/gem
#!/bin/sh

set -eu
sleep $(( $RANDOM % 7200 + 7200 ))

if type gem > /dev/null; then
  gem update --system
  gem update
  gem cleanup
fi
