# apkovl/etc/periodic/daily/gem
#!/bin/sh

set -eu
sleep 1200

if type gem > /dev/null; then
  gem update --system
  gem update
  gem cleanup
fi
