# apkovl/etc/periodic/daily/npm
#!/bin/sh

set -eu
sleep 2400

if type npm > /dev/null; then
  # npm install -g npm # updated by apk
  npm update -g
fi
